# macos/com.tradehaltalerts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.tradehaltalerts</string>
    <key>ProgramArguments</key>
    <array>
      <string>/usr/bin/python3</string>
      <string>/Users/jelk/trade-halt-alerts/scripts/halt_alerts.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/jelk/trade-halt-alerts</string>
    <key>EnvironmentVariables</key>
    <dict>
      <key>FMP_API_KEY</key>
      <string></string>
      <key>ALPHAVANTAGE_API_KEY</key>
      <string>RGQTSJ96HIHF178F</string>
      <key>TWELVEDATA_API_KEY</key>
      <string>d4d73d484f0b4fc88bea0a8b151e0fd9</string>
      <key>SEC_USER_AGENT</key>
      <string>TradeHaltAlerts/1.0 (jelkaim@gmail.com)</string>
      <key>X_API_BEARER_TOKEN</key>
      <string></string>
      <key>OPENAI_API_KEY</key>
      <string></string>
      <key>HALT_ALERTS_AI_CATALYST</key>
      <string>0</string>
      <key>TERMINAL_NOTIFIER_PATH</key>
      <string></string>
      <key>TERMINAL_NOTIFIER_SENDER</key>
      <string></string>
      <key>TERMINAL_NOTIFIER_ACTIVATE</key>
      <string></string>
      <key>HALT_ALERTS_OPEN_DETAILS</key>
      <string></string>
      <key>HALT_ALERTS_DETAILS_PORT</key>
      <string>8787</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jelk/trade-halt-alerts/logs/launchagent.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/jelk/trade-halt-alerts/logs/launchagent.err</string>
  </dict>
</plist>
